# etc/supervisor/me.easytee.service.pic_merge_test.conf
[program:me.easytee.service.pic_merge_test]
command=php /data/projects/easytee/easytee_v2/bin/cli.php --model=service/pic --action=run --env=test
numprocs=1
directory=/data/projects/easytee/easytee_v2/bin/
autostart=true
autorestart=true
startretries=10
exitcodes=0
stopsignal=KILL
stopwaitsecs=10
redirect_stderr=true